# deb-configs/postinst.in
#!/bin/bash


SERVICE_NAME="com.example.CalculatorService"
UNIT_PATH="/etc/systemd/system/calculator.service"
POLICY_PATH="/etc/dbus-1/system.d/com.example.CalculatorService.conf"


mkdir -p /etc/systemd/system/
cat << EOF > /etc/systemd/system/@PROJECT_NAME@.service
[Unit]
Description=@PROJECT_DESCRIPTION@
After=dbus.service
Requires=dbus.service

[Service]
Type=dbus
User=@PROJECT_NAME@
Group=@PROJECT_NAME@
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@
Restart=always
RestartSec=3s
WorkingDirectory=@CMAKE_INSTALL_PREFIX@/bin
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
[Install]
WantedBy=multi-user.target
EOF

systemctl daemon-reload

echo "Добавляем D-Bus policy"
mkdir -p /etc/dbus-1/system.d
cat > "$POLICY_PATH" <<EOF
<!DOCTYPE busconfig PUBLIC "-//freedesktop//DTD D-Bus Bus Configuration 1.0//EN"
 "http://www.freedesktop.org/standards/dbus/1.0/busconfig.dtd">
<busconfig>

  <policy user="calculator">
    <allow own="${SERVICE_NAME}"/>
  </policy>

  <policy context="default">
    <allow send_destination="${SERVICE_NAME}"/>
  </policy>

</busconfig>
EOF

echo "Создаём пользователя и группу calculator (если ещё не существуют)"
if ! getent group calculator >/dev/null; then
  groupadd -f  @PROJECT_NAME@
fi

if ! id -u calculator >/dev/null 2>&1; then
  useradd -m -k /dev/null -g @PROJECT_NAME@ @PROJECT_NAME@
fi



echo 'Устанавливаем права на файлы'
chown -R @PROJECT_NAME@:@PROJECT_NAME@ @CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@
chmod 755 @CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@

echo "Перечитываем конфигурации systemd"
systemctl daemon-reload

echo "Перезапускаем системный D-Bus для применения policy"
systemctl restart dbus.service

echo 'Запускаем @PROJECT_NAME@'
systemctl enable @PROJECT_NAME@
systemctl start @PROJECT_NAME@

sleep 1

echo 'Состояние сервиса @PROJECT_NAME@'
systemctl --no-pager status @PROJECT_NAME@