# deb-configs/prerm.in
#!/bin/bash
echo 'Останавливаем @PROJECT_NAME@'
if systemctl is-active --quiet @PROJECT_NAME@; then
        systemctl stop @PROJECT_NAME@
fi

# Also disable the service
systemctl disable @PROJECT_NAME@