# deb-configs/prerm
#!/bin/bash
echo 'Останавливаем calculator'
if systemctl is-active --quiet calculator; then
        systemctl stop calculator
fi

# Also disable the service
systemctl disable calculator